# systemd/odoo-semantic-webui.service.template
[Unit]
Description=Odoo Semantic Web UI (Admin Interface)
After=network.target odoo-semantic-mcp.service
Wants=odoo-semantic-mcp.service

[Service]
Type=simple
User=%i
WorkingDirectory=%h/git/odoo-semantic-mcp
EnvironmentFile=-%h/.odoo-semantic/.env
ExecStart=%h/.venv/odoo-semantic-mcp/bin/python -m src.web_ui
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
