# systemd/odoo-semantic-mcp.service.template
[Unit]
Description=Odoo Semantic MCP Server
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=%i
WorkingDirectory=%h/git/odoo-semantic-mcp
EnvironmentFile=-%h/.odoo-semantic/.env
ExecStart=%h/.venv/odoo-semantic-mcp/bin/python -m src.mcp.server
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
